# infra/host/systemd/lynis-audit.service
[Unit]
Description=Weekly Lynis security audit
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/lynis audit system --quiet --no-colors --logfile /var/log/lynis-weekly.log
Nice=15
IOSchedulingClass=idle
